# fs-skeleton/etc/network/if-up.d/ipv4ll
#!/bin/sh

if [ x${IF_ZCIP_PREF} != x ]
then
        /sbin/zcip -q -r ${IF_ZCIP_PREF} ${IFACE} /sbin/zcip.script
fi
