# fs-skeleton/etc/network/if-up.d/ntpd
#!/bin/sh

[ "$IFACE" != "lo" ] || exit 0

/usr/sbin/ntpd -q -p pool.ntp.org
