# fs-skeleton/etc/network/if-up.d/ntpd
#!/bin/sh

/usr/sbin/ntpd -q -p pool.ntp.org
